# dist/etc/systemd/system/framebufferserver.service
[Unit]
Description=reMarkable 2 Framebuffer Server
Before=xochitl.service launcher.service remarkable-reboot.service remarkable-shutdown.service
After=home.mount
StartLimitInterval=30
StartLimitBurst=5
ConditionFileNotEmpty=/home/crypto/lib/librm2fb_server.so

[Service]
Type=notify
ExecStart=/usr/bin/xochitl
Restart=on-failure
RestartSec=5
Environment="HOME=/home/root"
Environment="LD_PRELOAD=/home/crypto/lib/librm2fb_server.so"

[Install]
WantedBy=multi-user.target